# can-keyboard.service
[Unit]
Description=RNS-E CAN-Bus Keyboard Simulation (uinput)
After=can-handler.service
Wants=can-handler.service

[Service]
ExecStart=/usr/bin/python3 ${REPO_DIR}/can_keyboard.py
WorkingDirectory=${REPO_DIR}
Restart=always
RestartSec=3
User=pi
Group=input

[Install]
WantedBy=multi-user.target
